# PX4_DXP/bag-autorecord.service
[Unit]
Description=Auto rosbag recorder (records every mission start→complete)
# Needs the API (rover-server) for mission status, and MAVROS/RPP topics (px4-dxp).
After=rover-server.service px4-dxp.service
Wants=rover-server.service
PartOf=px4-dxp.service

[Service]
Type=exec
User=flash
Group=flash
WorkingDirectory=/home/flash/PX4_DXP
ExecStart=/home/flash/PX4_DXP/tools/bag_autorecord.sh
Restart=on-failure
RestartSec=5
StartLimitInterval=300
StartLimitBurst=10

# Security hardening (matches rpp-pipeline). Home is read-only EXCEPT the bag dir.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/flash/.ros /home/flash/bags_jet /home/flash/PX4_DXP /tmp /var/tmp
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Resource limits (recording is I/O, not CPU-critical — no realtime needed)
LimitNOFILE=65536
MemoryMax=1G
CPUQuota=100%

# Environment
Environment=ROS_DOMAIN_ID=0
Environment=ROS_LOCALHOST_ONLY=0
Environment=FASTRTPS_DEFAULT_PROFILES_FILE=/home/flash/PX4_DXP/config/fastdds_no_shm.xml
Environment=HOME=/home/flash
Environment=BAGS_DIR=/home/flash/bags_jet
Environment=PX4_DXP_DIR=/home/flash/PX4_DXP
Environment=MISSION_DEBUG_CONTROL_DIR=/home/flash/PX4_DXP/runtime/mission-debug
Environment=BAG_QOS_OVERRIDES=/home/flash/PX4_DXP/config/rosbag_qos_overrides.yaml
Environment=ROVER_TOKEN_FILE=/home/flash/PX4_DXP/config/rover_token
Environment=BAG_POST_TAIL_S=2.0
# Storage limits are site-configurable via a systemd drop-in.
Environment=BAG_MIN_FREE_BYTES=5368709120
Environment=BAG_LOW_FREE_BYTES=2147483648
Environment=BAG_MAX_TOTAL_BYTES=53687091200

# Output
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bag-autorecord

# Shutdown — give the daemon time to SIGINT-finalise an in-progress bag
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
