# PX4_DXP/rover-server.service
[Unit]
Description=Drawing Rover FastAPI Server (REST + Socket.IO + ROS2 bridge)
After=rpp-pipeline.service
Wants=rpp-pipeline.service
# Server can survive RPP restart (reconnects via TRANSIENT_LOCAL topics),
# so we do NOT use PartOf here. Only stop the server explicitly.

[Service]
Type=simple
User=flash
Group=flash
WorkingDirectory=/home/flash/PX4_DXP/server
# NOTE: rover-server is intentionally NOT gated on /mavros discovery. It is the
# operator's UI/telemetry surface and must come up even when MAVROS is down so
# /api/mission/status can report fcu_connected:False (graceful degradation).
# Upholds= (from px4-dxp) + Restart= reconnect once MAVROS returns.
ExecStart=/home/flash/PX4_DXP/server/run.sh
Restart=on-failure
RestartSec=5
StartLimitInterval=300
StartLimitBurst=10


# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/flash/.ros /home/flash/PX4_DXP /tmp /var/tmp
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=no

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=2G
CPUQuota=200%

# Environment
Environment=ROS_DOMAIN_ID=0
Environment=ROS_LOCALHOST_ONLY=0
Environment=FASTRTPS_DEFAULT_PROFILES_FILE=/home/flash/PX4_DXP/config/fastdds_no_shm.xml
Environment=HOME=/home/flash
Environment=FASTAPI_PORT=5001
Environment=ROVER_TOKEN_FILE=/home/flash/PX4_DXP/config/rover_token
Environment=MISSION_CAPTURE_REQUIRED=1
Environment=MISSION_DEBUG_CONTROL_DIR=/home/flash/PX4_DXP/runtime/mission-debug
Environment=MISSION_CAPTURE_ACK_TIMEOUT_S=3.0
# Auth disabled during dev bring-up. Remove this line to enforce X-Rover-Token
# (token auto-generated at the ROVER_TOKEN_FILE path above on first boot).
Environment=ROVER_DISABLE_AUTH=1

# Output
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rover-server

# Shutdown
KillMode=mixed
KillSignal=SIGTERM
TimeoutStartSec=30
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
