# PX4_DXP/rpp-pipeline.service
[Unit]
Description=RPP Controller Pipeline (twist_to_setpoint + rpp_controller + xtrack_logger)
After=px4-dxp.service
Wants=px4-dxp.service
# PartOf: if px4-dxp is stopped/restarted, this service stops too.
# RPP nodes are useless without MAVROS topics.
PartOf=px4-dxp.service

[Service]
Type=exec
User=flash
Group=flash
WorkingDirectory=/home/flash/PX4_DXP
ExecStartPre=/bin/bash -lc 'source /opt/ros/humble/setup.bash && timeout 35 python3 /home/flash/PX4_DXP/tools/ros2_mavros_health.py --timeout 30 node /mavros'
ExecStart=/home/flash/PX4_DXP/rpp_start.sh
Restart=on-failure
RestartSec=5
StartLimitInterval=300
StartLimitBurst=10

# Security hardening (same as px4-dxp)
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/flash/.ros /home/flash/PX4_DXP /tmp /var/tmp
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=no

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=2G
CPUQuota=200%

# P2.5 — Real-time scheduling on Jetson Orin
# Reduces timer jitter from ±20 ms to ±2 ms, enabling tighter control loops.
# FIFO priority 80 is high but not system-critical (kernel threads run at 50-60).
# Core 4 is a performance core on Jetson Orin (cores 0-3 are efficiency cores).
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=80
CPUAffinity=4

# Environment
Environment=ROS_DOMAIN_ID=0
Environment=ROS_LOCALHOST_ONLY=0
Environment=FASTRTPS_DEFAULT_PROFILES_FILE=/home/flash/PX4_DXP/config/fastdds_no_shm.xml
Environment=HOME=/home/flash

# Output
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rpp-pipeline

# Shutdown
KillMode=mixed
KillSignal=SIGTERM
TimeoutStartSec=45
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
